# scripts/rpi_controller.service
[Unit]
Description=RPI Controller: Binary to control the Zumo
After=syslog.target
After=rsyslog.service
After=network.target
After=rpi_updater.service

[Service]
Type=simple
WorkingDirectory=/
ExecStart=/root/local/main.out
Restart=always               

